# deploy/systemd/user/trafficpulse-processing.service
[Unit]
Description=TrafficPulse processing pipeline (aggregate/materialize/baselines)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/web-projects/traffic-pulse
ExecStartPre=/usr/bin/bash -lc 'until /usr/bin/docker info >/dev/null 2>&1; do sleep 5; done'
ExecStart=/usr/bin/bash -lc '\
  set +e; \
  /usr/bin/docker compose run --rm --no-deps --build api python scripts/ingest_sources_all.py || true; \
  /usr/bin/docker compose run --rm --no-deps --build api python scripts/aggregate_observations.py || true; \
  /usr/bin/docker compose run --rm --no-deps --build api python scripts/materialize_defaults.py || true; \
  /usr/bin/docker compose run --rm --no-deps --build api python scripts/build_baselines.py || true; \
  /usr/bin/docker compose run --rm --no-deps --build api python scripts/build_segment_quality.py || true; \
'
