# deploy/systemd/user/trafficpulse-monitor.service
[Unit]
Description=TrafficPulse health monitor (writes data/cache/alerts.log)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/web-projects/traffic-pulse
ExecStartPre=/usr/bin/bash -lc 'until /usr/bin/docker info >/dev/null 2>&1; do sleep 5; done'
ExecStart=/usr/bin/docker compose run --rm --no-deps --build api python scripts/monitor_health.py
