# deploy/systemd/user/trafficpulse-api.service
[Unit]
Description=TrafficPulse API (Docker Compose)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/web-projects/traffic-pulse
ExecStartPre=/usr/bin/bash -lc 'until /usr/bin/docker info >/dev/null 2>&1; do sleep 5; done'
ExecStart=/usr/bin/docker compose up -d api
ExecStop=/usr/bin/docker compose stop api
RemainAfterExit=yes
Restart=on-failure
RestartSec=15

[Install]
WantedBy=default.target

