# deploy/systemd/user/trafficpulse-monitor.timer
[Unit]
Description=TrafficPulse health monitor timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true
Unit=trafficpulse-monitor.service

[Install]
WantedBy=timers.target

